# container/el7-runner-node/systemd/var-cache-mock.mount
[Unit]
Description=Mount mock cache dir
Before=jenkins-jnlp-agent.service
Requires=mk-shared-slot.service
After=mk-shared-slot.service

[Mount]
What=/slt/var/cache/mock
Where=/var/cache/mock
Options=bind

[Install]
WantedBy=multi-user.target
RequiredBy=jenkins-jnlp-agent.service
